# scripts/llama-swap-mode-mcp.service
[Unit]
Description=llama-swap mode-switch MCP server (HTTP) — switch serving modes from MCP clients
Documentation=https://github.com/bradrlaw/ai-server
# Talks to llama-swap on 127.0.0.1:9090 and edits config/llama-swap.yaml on the
# host. Order after llama-swap so warms succeed at boot; it retries per-call anyway.
After=network-online.target llama-swap.service
Wants=network-online.target

[Service]
Type=simple
User=brad
Group=brad
WorkingDirectory=/srv/ai
# Serve over streamable-http (default is stdio). mcpo proxies this to Open WebUI
# and Copilot BYOK clients can register it over HTTP.
Environment=MODE_MCP_TRANSPORT=streamable-http
# Bind on all interfaces for LAN/Tailscale clients. NO auth (like comfyui-mcp /
# plan-build) — acceptable on the private LAN/Tailscale only.
Environment=MODE_MCP_HOST=0.0.0.0
Environment=MODE_MCP_PORT=9120
# The switcher (llama-swap-mode.py) needs PyYAML; run it with the system python3
# which has it (the mcp venv may not).
Environment=MODE_SWITCH_PY=/usr/bin/python3
Environment=LLAMASWAP_URL=http://127.0.0.1:9090
Environment=HOME=/home/brad
# Reuse the comfyui-mcp venv (already has the `mcp` package); no new venv.
ExecStart=/srv/ai/venvs/comfyui-mcp/bin/python /srv/ai/docker/mcpo/llama_swap_mode_mcp.py
Restart=on-failure
RestartSec=3
TimeoutStopSec=15
KillMode=mixed

[Install]
WantedBy=multi-user.target
